# supervisord.conf
[supervisord]
nodaemon=true

[program:mysqlproxy]
command=mysql-proxy
