# supervisord.conf
[supervisord]
nodaemon=true

[program:mysqlproxy]
command=mysql-proxy --defaults-file=/etc/mysql-proxy.conf --daemon=false
